# cloudquery/crontab
0 * * * * /app/sync.sh /app/hourly.yml
0 9 * * * /app/sync.sh /app/daily.yml
0 18 * * * /app/sync.sh /app/daily.yml
* * * * * /app/sync.sh /app/frequent.yml
