# cloudquery/crontab
0 * * * * /app/sync.sh /app/hourly.yml
0 9 * * * /app/sync.sh /app/daily.yml
